# pi/config/micro-bot.service
[Unit]
Description=MICRO BOT Bridge Server - Swarm Communication Hub
Documentation=file:///home/pi/micro_bot/INFO.txt
After=network.target hostapd.service dnsmasq.service
Wants=hostapd.service dnsmasq.service
Requires=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/micro_bot
Environment=PATH=/home/pi/micro_bot/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/pi/micro_bot
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/micro_bot/venv/bin/python bridge_server.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
Restart=always
RestartSec=10
TimeoutStopSec=30
TimeoutStartSec=60

# Logging configuration
StandardOutput=journal
StandardError=journal
SyslogIdentifier=micro-bot

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=512M
CPUQuota=200%

# Health monitoring
WatchdogSec=30

[Install]
WantedBy=multi-user.target